# install/files/crontab
*/5 * * * * absurdo python3 /home/absurdo/apache/htdocs/pastafari/modules/pastafari/cron/check_down.py
