# install/files/crontab
*/5 * * * * pzoo python3 /path/to/pastafari/modules/pastafari/cron/check_down.py
