# scripts/monit/debian_jessie/files/crontab/alive
*/2 * * * * pzoo /home/spanel/modules/pastafari/scripts/monit/debian_jessie/files/get_info.py

